# crontab
* * * * * wget -qO- $RESPONSE_URL/slack/cron_minute